# codedir/environments/production/modules/libs_data/files/users/silence.service
[Unit]
Description=Play Silence

[Service]
ExecStart=%h/.config/systemd/user/silence.sh

[Install]
WantedBy=default.target
